# etc_init_webhookr.conf
description "statsd"

start on startup
stop on shutdown

respawn

script
    chdir /usr/local/webhookr
    exec python manage.py runserver_socketio
end script
EOF